# config/relampagos-backend.service
[Unit]
Description=Relampagos Backend API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/html/relampagos/backend
ExecStart=/usr/bin/bun index.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=PATH=/usr/bin:/usr/local/bin:/home/www-data/.bun/bin

[Install]
WantedBy=multi-user.target
